# omero_supervisor.conf
[supervisord]
nodaemon=true

[program:omero]
command=/home/omero/OMERO.server/bin/omero admin start
user=omero

[program:omero_web]
command=/home/omero/OMERO.server/bin/omero web start
user=omero

